# 03 自动化运维组件/ansible/roles/01 roles/mongodb/templates/mongodb.service.j2
[Unit]
Description=start mongo cluster
After=network.target

[Service]
ExecStart={{ var_runscript }}/run_mongoCluster.sh start
ExecStop={{ var_runscript }}/run_mongoCluster.sh stop
Type=forking

[Install]
WantedBy=multi-user.target